# templates/node_exporter.service.j2
[Unit]
Description=Node Exporter
After=network.target

[Service]
User=root
Group=root
Type=simple
ExecStart=/usr/local/bin/node_exporter {{ service_parameters }}

[Install]
WantedBy=multi-user.target
